# scripts/pve/lxc-clone-queue/lxc-clone-queue.service
[Unit]
Description=LXC Clone Queue - PVE clone request serialization proxy
Documentation=https://github.com/manaflow-ai/cmux
After=network-online.target pve-cluster.service pvedaemon.service
Wants=network-online.target
Requires=pvedaemon.service

[Service]
Type=simple
ExecStart=/usr/local/bin/lxc-clone-queue \
    -listen :8081 \
    -pve-addr https://127.0.0.1:8006 \
    -max-retries 5 \
    -task-poll-interval 2s \
    -task-timeout 5m \
    -request-timeout 10m \
    -insecure-tls=true

Restart=on-failure
RestartSec=5
TimeoutStartSec=10
TimeoutStopSec=10

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lxc-clone-queue

# Environment (can be overridden in /etc/default/lxc-clone-queue)
EnvironmentFile=-/etc/default/lxc-clone-queue

[Install]
WantedBy=multi-user.target
