# configs/profile.d/cmux-env.sh
# cmux sandbox environment variables
# Auto-loaded on shell startup via /etc/profile.d/
export IS_SANDBOX=1
