# configs/systemd/cmux-xterm.service
[Unit]
Description=Cmux XTerm Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/opt/app/app.env
LogsDirectory=cmux
ExecStartPre=/bin/mkdir -p /var/log/cmux
ExecStart=/bin/bash -lc 'exec /usr/local/bin/cmux-xterm'
Restart=on-failure
RestartSec=2
StandardOutput=append:/var/log/cmux/cmux-xterm.log
StandardError=append:/var/log/cmux/cmux-xterm.log

[Install]
WantedBy=cmux.target