# configs/systemd/cmux-token-generator.service
[Unit]
Description=CMUX Auth Token Generator
Before=cmux-execd.service cmux-worker.service cmux-worker-daemon.service
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/cmux-token-init

[Install]
WantedBy=multi-user.target
